# headnodes/hpc-meowth1/ansible/roles/slurm_db/files/slurmdbd.service.j2
[Unit]
Description=slurmdbd
After=network.target
[Service]
Type=forking
User=root
ExecStart=/opt/slurm/sbin/slurmdbd
Restart = always
[Install]
WantedBy=multi-user.target
